# conf/app.supervisor.conf
[program:watch_corona]
command=/home/ubuntu/watch-corona/server/run.sh
autostart=true
autorestart=true
stdout_logfile=/home/ubuntu/logs/watchcorona.log
redirect_stderr=true
