# conf/app.supervisor.conf
[program:watch_corona]
command=/home/rz1993/watch-corona/server/run.sh
autostart=true
autorestart=true
stdout_logfile=/home/rz1993/logs/watchcorona.log
redirect_stderr=true
